# etc/cron.d/cashrewards_update
# Update the Cashrewards blocklist
0 5 * * * pihole /home/pi/adlist_converter.py
